# pi-api/emg-api.service
[Unit]
Description=EMG Arduino API — Serial to Neon DB (Python)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ahmed
WorkingDirectory=/home/ahmed/muscle-sensor-ui/pi-api
ExecStart=/usr/bin/python3 /home/ahmed/muscle-sensor-ui/pi-api/server.py
Restart=always
RestartSec=5
ExecStartPre=/bin/sleep 5
EnvironmentFile=/home/ahmed/muscle-sensor-ui/pi-api/.env
StandardOutput=append:/var/log/emg-api.log
StandardError=append:/var/log/emg-api.log
SyslogIdentifier=emg-api

[Install]
WantedBy=multi-user.target
